# deploy/g2web.service
[Unit]
Description=G2 Web-Editor Bridge
After=network.target

[Service]
Type=simple
User=g2web
Group=plugdev
# JVM-Properties müssen via JAVA_OPTS rein (Gradle-Startskript reicht Programm-Args nicht als -D durch)
Environment=JAVA_OPTS=-Dg2web.service=usb -Dg2web.port=8080
# Falls der Vorgänger unsauber starb, hängt der G2 sonst in einem nicht-verbindbaren Zustand
ExecStartPre=+/bin/sh -c '/usr/bin/usbreset 0ffc:0002 || true'
ExecStart=/opt/g2web/bin/g2web-backend
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
